# systemd/litellm.service
[Unit]
Description=LiteLLM proxy for pawn-agent
After=network.target pawn-agent.service

[Service]
Type=simple
User=%u
WorkingDirectory={{PROJECT_DIR}}
Environment=PAWN_AGENT_URL=http://localhost:8000
ExecStart={{BIN_DIR}}/litellm --config {{PROJECT_DIR}}/litellm_config.yaml --port 4000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
